# systemd/practice-bank-autodebit-retry.service
[Unit]
Description=PRACTICE BANK Monthly Auto-debit Retry Run (= 15-AD retry-only mode)
After=network.target postgresql.service

[Service]
Type=oneshot
User=batch
Group=batch
WorkingDirectory=/opt/practice-bank
EnvironmentFile=/opt/practice-bank/etc/practice-bank.env
ExecStart=/usr/bin/make -C /opt/practice-bank/subsystems/22-operations master-load-customers
StandardOutput=journal
StandardError=journal
TimeoutStartSec=2h

[Install]
WantedBy=multi-user.target
