# systemd/practice-bank-batch-daily.service
[Unit]
Description=PRACTICE BANK Daily Batch Run (= ops-batch-daily pipeline)
After=network.target postgresql.service rabbitmq-server.service
Wants=postgresql.service

[Service]
Type=oneshot
User=batch
Group=batch
WorkingDirectory=/opt/practice-bank
EnvironmentFile=/opt/practice-bank/etc/practice-bank.env
ExecStart=/opt/practice-bank/subsystems/22-operations/bin/run-batch-daily-wrapper.sh daily
StandardOutput=journal
StandardError=journal
TimeoutStartSec=4h
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
